# roam2/wpa_supplicant.service
[Unit]
Description=WPA Supplicant for WiFi Authentication
After=network.target
Before=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/opt/roam/wpa_supplicant_wrapper.sh
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=10
User=root
Group=root
StandardOutput=journal
StandardError=journal
SyslogIdentifier=wpa_supplicant

# Relaxed security settings for wpa_supplicant compatibility
NoNewPrivileges=false
PrivateTmp=false
ProtectSystem=false
ProtectHome=false

# Environment
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin

[Install]
WantedBy=multi-user.target
